# autooverview.service
[Unit]
Description=AutoOverview Backend Service (FastAPI + uvicorn)
Documentation=https://github.com/zhancongc/AutoOverview
After=network-online.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/app/AutoOverview/backend

# Environment
Environment="PATH=/app/AutoOverview/backend/.venv/bin:/usr/local/bin:/usr/bin"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=-/app/AutoOverview/backend/.env

# Start command - uvicorn in production mode (no reload)
# Worker 数量建议: (2 × CPU cores) + 1
# - 2核 VPS: --workers 4
# - 4核 VPS: --workers 6  ← 当前配置（预留资源给其他服务）
# - 8核 VPS: --workers 16
ExecStart=/app/AutoOverview/backend/.venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000 --workers 6

# Restart policy
Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitInterval=60

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
TimeoutStartSec=60

# Logging
StandardOutput=append:/app/AutoOverview/backend/logs/uvicorn_access.log
StandardError=append:/app/AutoOverview/backend/logs/uvicorn_error.log

# Security
PrivateTmp=true

[Install]
WantedBy=multi-user.target
